# docs/launchd/red-crawler.collect-night.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.red-crawler.collect-night</string>

  <key>WorkingDirectory</key>
  <string>/Users/tommy/Documents/GitHubOpenSources/red-crawler</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>python</string>
    <string>-m</string>
    <string>red_crawler.cli</string>
    <string>collect-nightly</string>
    <string>--storage-state</string>
    <string>/Users/tommy/Documents/GitHubOpenSources/red-crawler/state.json</string>
    <string>--db-path</string>
    <string>/Users/tommy/Documents/GitHubOpenSources/red-crawler/data/red_crawler.db</string>
    <string>--report-dir</string>
    <string>/Users/tommy/Documents/GitHubOpenSources/red-crawler/reports</string>
    <string>--cache-dir</string>
    <string>/Users/tommy/Documents/GitHubOpenSources/red-crawler/.cache/red-crawler</string>
    <string>--crawl-budget</string>
    <string>4</string>
    <string>--daily-account-budget</string>
    <string>12</string>
    <string>--daily-search-term-budget</string>
    <string>2</string>
    <string>--startup-jitter-minutes</string>
    <string>90</string>
    <string>--slot-name</string>
    <string>night</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/tommy/Documents/GitHubOpenSources/red-crawler/logs/collect-night.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/tommy/Documents/GitHubOpenSources/red-crawler/logs/collect-night.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
